# apps/subscriber/config/pull_subscriber.ini
[program:pull_subscriber]
command=python apps/subscriber/pull/pull_subscriber.py
directory=/home/hogehoge/releases/current/apps/
environment=PYTHONPATH=/home/hogehoge/releases/current/apps/
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/pull_subscriber/pull_subscriber.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true
